# roles/ui/files/etc/skel/.config/systemd/user/guake.service
[Unit]
Description="start the guake terminal"

[Service]
ExecStart=/usr/bin/guake

[Install]
WantedBy=sway-session.target
