# roles/ui/files/etc/skel/.config/systemd/user/radiotray-ng.service
[Unit]
Description="start the radiotray"

[Service]
ExecStart=/usr/bin/radiotray-ng

[Install]
WantedBy=sway-session.target
